# deploy/dockercmd.service
[Unit]
Description=Docker Commander — Docker monitoring & control panel
Documentation=https://github.com/koduj-dev/docker-commander
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
# Dedicated unprivileged user that is a member of the `docker` group so it can
# reach the daemon socket. Create it once with:
#   sudo useradd --system --no-create-home --shell /usr/sbin/nologin dockercmd
#   sudo usermod -aG docker dockercmd
User=dockercmd
Group=dockercmd
SupplementaryGroups=docker

# The -data-dir flag is pinned here because it must match ReadWritePaths below.
# Everything else (DC_ADDR, DC_REDIS_ADDR, DC_METRICS_TOKEN, …) can come from a
# config file or the environment file. The binary reads
# /etc/docker-commander/commander.conf by default (override with -config / the
# DC_CONFIG env var). Precedence: flags > environment > config file.
ExecStart=/usr/local/bin/dockercmd -data-dir /var/lib/dockercmd
EnvironmentFile=-/etc/dockercmd/dockercmd.env

Restart=on-failure
RestartSec=5s

# --- hardening ---------------------------------------------------------------
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/dockercmd
StateDirectory=dockercmd
# The daemon socket is reached via the docker group, not raw device access.

[Install]
WantedBy=multi-user.target
