# internal/service/unit.service
[Unit]
Description=Docker Commander — Docker monitoring & control panel
Documentation=https://github.com/koduj-dev/docker-commander
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
# Dedicated unprivileged user that is a member of the `docker` group so it can
# reach the daemon socket. Create it once with:
#   sudo useradd --system --no-create-home --shell /usr/sbin/nologin dockercmd
#   sudo usermod -aG docker dockercmd
User=dockercmd
Group=dockercmd
SupplementaryGroups=docker

# The -data-dir flag is pinned here because it must match ReadWritePaths below.
# Everything else (DC_ADDR, DC_REDIS_ADDR, DC_METRICS_TOKEN, …) is read from
# /etc/docker-commander/commander.conf — the binary loads it automatically
# (override the path with -config or DC_CONFIG). Precedence: flags > env > file.
ExecStart=/usr/local/bin/dockercmd -data-dir /var/lib/dockercmd

# Compose/Projects shell out to the `docker compose` CLI. The docker CLI bases
# its plugin discovery on its config dir (~/.docker by default), but ProtectHome
# below mounts this user's home (/home/dockercmd) inaccessible (EACCES) — and
# docker treats that as fatal, so `docker compose` looks like an unknown command
# and Deploy/Down get disabled. Point the config dir at the writable state dir so
# the compose plugin (in /usr/libexec/docker/cli-plugins) is found.
Environment=DOCKER_CONFIG=/var/lib/dockercmd/.docker

Restart=on-failure
RestartSec=5s

# Tag journal/syslog entries as "dockercmd" (journalctl -t dockercmd). Add
# ForwardToSyslog=yes in journald.conf to also push them to a syslog daemon.
SyslogIdentifier=dockercmd

# --- hardening ---------------------------------------------------------------
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/dockercmd
StateDirectory=dockercmd
# The daemon socket is reached via the docker group, not raw device access.

[Install]
WantedBy=multi-user.target
